# systemd/kroger-recipes-restart.path
[Path]
PathChanged=/opt/kroger-recipes/server

[Install]
WantedBy=multi-user.target